# docker/supervisord.only_server.conf
[supervisord]
nodaemon=true
user=root

[program:picocmt]
command=/app/server/picocmt --config-path /app/server/config.toml
autostart=true
autorestart=true
